# fgcore_runner/scripts/systemd/open5gs-udrd.service.in
[Unit]
Description=Open5GS UDR Daemon
Wants=mongodb.service mongod.service

[Service]
Type=simple

User=root
User=root

Restart=always
ExecStart=/bin/open5gs-udrd -c /home/ops/nf_configs/udr.yml
RestartSec=2
RestartPreventExitStatus=1
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
